# rc.local
# .rc.local
# Local customizations for interactive shells
# Eryn Wells <eryn@erynwells.me>


